# scripts/upstart/audiomasher.conf
description "AudioMasher"
author "Pierre Cusa"

start on (local-filesystems and net-device-up)
stop on runlevel [!2345]

respawn
chdir /srv/audiomasher/
exec /srv/audiomasher/AudioMasher